# conf/ddns.supervisor.conf
[program:ddns]
command=/usr/bin/node server.js
directory=/path/to/ddns
user=nsd
redirect_stderr=true
autostart=true
autorestart=true
startsecs=5
